# Aston/com.aston.daily-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aston.daily-rotate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/nikhilr5/Desktop/Kalshi2.0/Aston/daily_rotate.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/nikhilr5/Library/Logs/aston-daily-rotate.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nikhilr5/Library/Logs/aston-daily-rotate.log</string>

  <!-- aws CLI lives in /usr/local/bin or /opt/homebrew/bin; ensure
       both are on PATH so subprocess.run finds it from launchd. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/nikhilr5</string>
  </dict>
</dict>
</plist>
